# contrib/init/bitmonkeyd.conf
description "bitmonkey Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env BITMONKEYD_BIN="/usr/bin/bitmonkeyd"
env BITMONKEYD_USER="bitmonkey"
env BITMONKEYD_GROUP="bitmonkey"
env BITMONKEYD_PIDDIR="/var/run/bitmonkeyd"
# upstart can't handle variables constructed with other variables
env BITMONKEYD_PIDFILE="/var/run/bitmonkeyd/bitmonkeyd.pid"
env BITMONKEYD_CONFIGFILE="/etc/bitmonkey/bitmonkey.conf"
env BITMONKEYD_DATADIR="/var/lib/bitmonkeyd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # bitmonkeyd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$BITMONKEYD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run bitmonkeyd."
        echo "The setting must appear in $BITMONKEYD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"bitmonkey Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$BITMONKEYD_PIDDIR"
    chmod 0755 "$BITMONKEYD_PIDDIR"
    chown $BITMONKEYD_USER:$BITMONKEYD_GROUP "$BITMONKEYD_PIDDIR"
    chown $BITMONKEYD_USER:$BITMONKEYD_GROUP "$BITMONKEYD_CONFIGFILE"
    chmod 0660 "$BITMONKEYD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$BITMONKEYD_PIDFILE" \
    --chuid $BITMONKEYD_USER:$BITMONKEYD_GROUP \
    --exec "$BITMONKEYD_BIN" \
    -- \
    -pid="$BITMONKEYD_PIDFILE" \
    -conf="$BITMONKEYD_CONFIGFILE" \
    -datadir="$BITMONKEYD_DATADIR" \
    -disablewallet \
    -daemon

